# backup/crontab
0 3 * * * /app/backup/backup_db.sh >> /var/log/cron.log 2>&1
